# deploy/systemd/vegan-basket-ingest.timer
[Unit]
Description=Run Vegan Basket Analytics ingest daily at 23:00 IST

[Timer]
# Fire at 23:00 IST. IST is UTC+05:30 with no DST, so use explicit UTC here
# (matches Dockerfile cron and the crontab example below).
OnCalendar=*-*-* 17:30:00 UTC

# If the machine was off at fire time, run as soon as it comes back up.
Persistent=true

[Install]
WantedBy=timers.target
